# misc/miredo.in
#! /bin/sh
# $Id$
#
# miredo start/stop script for Debian GNU/Linux
# Author: Rémi Denis-Courmont <rdenis (at) simphalempin (dot) com>
#
# chkconfig: 345 17 83
# description: Starts and stops the Miredo daemon \
#	       used to provide IPv6 tunneling over UDP through NAT.
#
### BEGIN INIT INFO
# Provides:          teredo ipv6
# Requires:          $local_fs $network $syslog $time
# Short-Description: Teredo IPv6 tunnel
# Description:       Miredo daemon for tunneling of IPv6 through NAT
#                    within UDP/IPv4, as specified by the Teredo mechanism.
# Default-Start:     3 4 5
# Default-Stop:      S 0 1 2 6
### END INIT INFO


PATH=/usr/sbin:/usr/bin:/sbin:/bin
DESC="Teredo IPv6 tunneling daemon"
NAME=miredo
DAEMON=@sbindir@/$NAME
DAEMON_ARGS=""
PIDFILE=@localstatedir@/run/$NAME.pid
SCRIPTNAME=/etc/init.d/$NAME

[ -x "$DAEMON" ] || exit 0

# Source defaults.
[ -r /etc/default/$NAME ] && . /etc/default/$NAME

# Ensure configuration is readable.
[ -r @sysconfdir@/$NAME.conf ] || exit 0

. /lib/lsb/init-functions

case "$1" in
  start)
	if [ "x$STARTMIREDO" != "xtrue" ]; then
		echo "STARTMIREDO is set to false in /etc/default/$NAME"
		echo "$DAEMON not starting"
		exit 0
	fi
	log_daemon_msg "Starting $DESC" "$NAME"
	start-stop-daemon --start --quiet --pidfile "$PIDFILE" \
		--exec "$DAEMON" --oknodo -- $DAEMON_ARGS
	log_end_msg $?
	;;
  stop)
	log_daemon_msg "Stopping $DESC" "$NAME"
	start-stop-daemon --stop --quiet --pidfile "$PIDFILE" \
		--retry 1 --oknodo
	log_end_msg $?
	;;
  reload|force-reload)
	log_daemon_msg "Reloading $DESC" "$NAME"
	start-stop-daemon --stop --signal 1 --quiet --pidfile "$PIDFILE" \
		--exec "$DAEMON"
	log_end_msg $?
	;;
  restart)
	$0 stop
	sleep 1
	$0 start
	;;
  *)
	echo "Usage: $SCRIPTNAME {start|stop|restart|reload|force-reload}" >&2
	exit 1
	;;
esac

exit $?

